# templates/init_d_hbase.erb
#!/bin/sh
#
# hbase - this script starts and stops the hbase-server daemon
#
# chkconfig:    - 85 15 
# description:  HBase is an open source, non-relational, distributed database modeled after Google's BigTable and is written in Java
# processname:  hbase-daemon
# config dir:   /usr/lib/hbase/conf
# pidfile:      /var/run/hbase.pid

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Source service specific environment variables
. /etc/profile.d/hbase.sh

# Check that networking is up.
[ "$NETWORKING" = "no" ] && exit 0

prog="hbase"
port="<%= @port %>"
hbase_root="/usr/lib/hbase"
zookeeper_root='/usr/lib/zookeeper'
hbase_daemon="$hbase_root/bin/hbase-daemon.sh"
zookeeper_daemon="$zookeeper_root/bin/zkServer.sh"
hbase_master_start="$hbase_daemon start master"
hbase_master_stop="$hbase_daemon stop master"
hbase_regionserver_start="$hbase_daemon start regionserver"
hbase_regionserver_stop="$hbase_daemon stop regionserver"
hbase_rest_start="$hbase_daemon start rest -p $port"
hbase_rest_stop="$hbase_daemon stop rest"
zookeeper_start="$zookeeper_daemon start"
zookeeper_stop="$zookeeper_daemon stop"

lockfile=/var/lock/subsys/hbase

start() {
    [[ -x $hbase_daemon ]] || exit 5
    echo -n $"Starting $prog: "
    su - zookeeper -c "$zookeeper_start"
    su - hbase -c "$hbase_master_start; $hbase_regionserver_start; $hbase_rest_start"
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
    return $retval
}

stop() {
    echo -n $"Stopping $prog: "
    su - zookeeper -c "$zookeeper_stop"
    su - hbase -c "$hbase_rest_stop; $hbase_regionserver_stop; $hbase_master_stop"
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}

restart() {
    stop
    start
}

status() {
    if [ -f $lockfile ]; then
      echo hbase is running
      exit 0
    else
      echo hbase is stopped
      exit 3
    fi
}

case "$1" in
    start)
        [[ -f $lockfile ]] && exit 0
        $1
        ;;
    stop)
        [[ -f $lockfile ]] || exit 0
        $1
        ;;
    restart)
        restart
        ;;
    status)
        status
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart}"
        exit 2
esac
